# systemd/teagram-mini-brain.service.in
# teagram-mini-brain — the Pipecat voice brain (STT->LLM->TTS gateway on :${BRAIN_PORT}/talk).
# Rendered by install.sh; ${VARS} come from @ETC@/brain.env. Secrets (the LLM key, the
# OpenClaw token) are read by the brain from ~/.config/teagram-mini at call time, never here.
[Unit]
Description=teagram-mini brain (voice pipeline gateway)
After=network-online.target teagram-mini-engine.service
Wants=network-online.target
Requires=teagram-mini-engine.service

[Service]
Type=simple
User=@USER@
EnvironmentFile=@ETC@/brain.env
ExecStart=@PREFIX@/venv/bin/teagram-mini-brain --port ${BRAIN_PORT}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
